# systemd/ezcams-pi-agent.service
[Unit]
Description=EZ Cams Raspberry Pi Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=EZCAMS_PI_CONFIG_DIR=/etc/ezcams-pi
ExecStart=/opt/ezcams-pi-agent/.venv/bin/ezcams-pi-agent run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
