# systemd/ezcams-pi-agent.service
[Unit]
Description=EZ Cams Raspberry Pi Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ezcams-pi-agent
ExecStartPre=/opt/ezcams-pi-agent/.venv/bin/ezcams-pi-agent ensure
ExecStart=/opt/ezcams-pi-agent/.venv/bin/ezcams-pi-agent run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
